# deploy/systemd/ai_papers.service
[Unit]
Description=AI Papers service
After=network.target

[Service]
User=www-data
WorkingDirectory=/srv/ai_papers
Environment="PATH=/srv/ai_papers/.venv/bin"
ExecStart=/srv/ai_papers/.venv/bin/gunicorn "deploy.wsgi:app" -w 4 -b 127.0.0.1:5000
Restart=on-failure

[Install]
WantedBy=multi-user.target
